# daemon/cleansing-fire-gatekeeper.service
# =============================================================================
# Cleansing Fire Gatekeeper Daemon — systemd service
# =============================================================================
# GPU access serializer for local Ollama. Listens on :7800.
#
# Install:
#   sudo cp daemon/cleansing-fire-gatekeeper.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable cleansing-fire-gatekeeper
#   sudo systemctl start cleansing-fire-gatekeeper
#
# Logs:
#   journalctl -u cleansing-fire-gatekeeper -f
# =============================================================================

[Unit]
Description=Cleansing Fire Gatekeeper Daemon (Ollama GPU serializer)
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=%i
WorkingDirectory=%h/cleansing-fire
ExecStart=/usr/bin/python3 daemon/gatekeeper.py --port 7800 --queue-size 5 --model mistral-large:123b
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=HOME=%h

[Install]
WantedBy=multi-user.target
